# test-data/diffs/elastic-logstash-742
diff --git a/pkg/logstash.sysv.redhat b/pkg/logstash.sysv.redhat
index cadbbd42a6a..80eeae7a0b3 100755
--- a/pkg/logstash.sysv.redhat
+++ b/pkg/logstash.sysv.redhat
@@ -4,7 +4,7 @@
 #
 #       Starts Logstash as a daemon
 #
-# chkconfig: 2345 20 80
+# chkconfig: 2345 90 10
 # description: Starts Logstash as a daemon
 
 ### BEGIN INIT INFO
